# gentoo/init.d/mogstored
#!/sbin/runscript

NAME="mogstored"
BINARY="/usr/bin/mogstored"

depend() {
	use net
	after mysql
	before mogilefsd
}

start() {
	ebegin "Starting $NAME"
	start-stop-daemon --start --quiet \
		--pidfile ${PIDFILE} \
		--exec ${BINARY} -- ${MOGSTORED_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping $NAME"
	start-stop-daemon --stop --quiet --pidfile="${PIDFILE}"
	eend $?
}

# vim: ft=gentoo-init-d:
